# etc/rc.local
/usr/bin/ss-server -c /etc/ss.conf &
/usr/sbin/tinyproxy  -c /etc/tinyproxy.conf &
/usr/bin/stunnel4 /etc/stunnel/stunnel.conf &
